# kde/.config/autostart/krunner-warmup.desktop
[Desktop Entry]
Exec=krunner -d
Name=Krunner Warm-up
Type=Application
X-KDE-AutostartScript=true
